# scheduling/com.ytmetrics.briefing.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Weekly EMAILED briefing — Mondays at 06:00 local (Eastern for RDU). Runs
  scripts/weekly_briefing.sh: pull → insights → briefing --email.

  One-time setup:
    1. Edit the uv path / WorkingDirectory below if they differ.
    2. Create a Gmail App Password (Google Account → Security → 2-Step Verification →
       App passwords) and store it WITHOUT committing it:
         echo 'your-app-password' > /Users/robert/src/ytmetrics/secrets/smtp_password
       (secrets/ is gitignored). Or export YTMETRICS_SMTP_PASSWORD in the job's env.
    3. Add an [email] section to config.toml (see config.example.toml).
    4. cp this file to ~/Library/LaunchAgents/com.ytmetrics.briefing.plist (drop ".example")
    5. launchctl load ~/Library/LaunchAgents/com.ytmetrics.briefing.plist
    6. launchctl start com.ytmetrics.briefing   # run once now to verify the email arrives

  Note: the Mac must be awake at 06:00 Monday; if asleep, launchd runs the job on the
  next wake. Needs a live OAuth token for the pull step (publish the consent screen to
  "In production" so the token doesn't expire weekly).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ytmetrics.briefing</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/robert/src/ytmetrics/scripts/weekly_briefing.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/robert/src/ytmetrics</string>

  <!-- Monday (Weekday 1) at 06:00 local time. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/robert/src/ytmetrics/logs/launchd.briefing.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/robert/src/ytmetrics/logs/launchd.briefing.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
